# systemd/log_cleaner.service
[Unit]
Description=Small script for directory cleaning with saving their structure.

[Service]
ExecStart=/opt/log_cleaner.sh
