# systemd/log_cleaner.timer
[Unit]
Description=Timer for running log_cleaner every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
